# deploy/systemd/quant-business-console.service
[Unit]
Description=MLBot read-only business console (Trade Map)
After=network.target

[Service]
Type=simple
User=quant
WorkingDirectory=/opt/quant-engine
Environment=PYTHONPATH=/opt/quant-engine/business-console/backend
Environment=MLBOT_CONSOLE_PORT=8800
Environment=MLBOT_CONSOLE_FEATURE_BUS_ROOT=/opt/quant-engine/live/shared_feature_bus
Environment=MLBOT_CONSOLE_LIVE_DATA_ROOT=/opt/quant-engine/live/highcap/data
Environment=MLBOT_CONSOLE_ENGINE_DATA_ROOT=/opt/quant-engine/data
Environment=MLBOT_CONSOLE_LIVE_ROOT=/opt/quant-engine/live/highcap
Environment=MLBOT_CONSOLE_GRAFANA_URL=http://127.0.0.1:3000
# Optional Basic Auth (uncomment):
# Environment=MLBOT_CONSOLE_BASIC_AUTH_USER=admin
# Environment=MLBOT_CONSOLE_BASIC_AUTH_PASSWORD=change-me
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 127.0.0.1 --port 8800
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
