# terraform/systemd/quant-engine.service
[Unit]
Description=Quant Engine - Three Strategy Live Trading (Docker)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
TimeoutStartSec=120

# 清理残留容器
ExecStartPre=-/usr/bin/docker rm -f quant-engine

# 运行容器
# - 镜像由 GitHub Actions 构建并推送，本地 tag 为 quant-engine:latest
# - /opt/quant-engine/live/highcap/data → 持久化运行时数据 (db, ticks, features)
# - /opt/quant-engine/data → funding_rate + open_interest parquet (研究数据)
# - /opt/quant-engine/live/binance_mainnet.env → API 密钥 (只读挂载)
ExecStart=/usr/bin/docker run \
    --name quant-engine \
    --network host \
    --memory=3g \
    --cpus=1.8 \
    -e MLBOT_LIVE_TRADE_SIZE=0.0 \
    -e MLBOT_ORDER_MODE=test \
    -e MLBOT_ORDER_MANAGER_ENABLED=true \
    -e PYTHONUNBUFFERED=1 \
    -v /opt/quant-engine/live/highcap/data:/app/live/highcap/data \
    -v /opt/quant-engine/data:/app/data \
    -v /opt/quant-engine/live/binance_mainnet.env:/app/live/binance_mainnet.env:ro \
    quant-engine:latest

# 停止容器
ExecStop=/usr/bin/docker stop -t 30 quant-engine
ExecStopPost=-/usr/bin/docker rm -f quant-engine

# 日志输出到 journald（供 Filebeat 采集）
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quant-engine

# 崩溃自动重启
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
