# etc/systemd/mlbot-monitor-notify@.service
[Unit]
Description=Telegram notify for failed mlbot monitor unit %i
After=network-online.target

[Service]
Type=oneshot
Environment=MLBOT_MONITORING_ENV=/opt/quant-engine/monitoring/.env
ExecStart=/opt/quant-engine/scripts/monitoring/monitor_telegram_notify.sh %i failed
TimeoutStartSec=30
