# etc/systemd/mlbot-weekly-watchdog.service
[Unit]
Description=mlbot weekly regime watchdog (alias of monitor schedule weekly)
After=docker.service network-online.target
Requires=docker.service
OnFailure=mlbot-monitor-notify@%n.service

[Service]
Type=oneshot
ExecStart=/opt/quant-engine/deploy/systemd/mlbot-monitor-docker-run.sh weekly

[Install]
WantedBy=multi-user.target
