# etc/systemd/mlbot-monitor-weekly.timer
[Unit]
Description=mlbot monitor weekly — Sunday 08:00 UTC

[Timer]
OnCalendar=Sun *-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
