# etc/systemd/mlbot-monitor-daily.timer
[Unit]
Description=mlbot monitor daily — 06:30 UTC

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
