# etc/systemd/mlbot-account-reconcile.timer
[Unit]
Description=MLBot ABC account reconciliation — every 30 minutes UTC

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
