# etc/systemd/mlbot-monitor-monthly.service
[Unit]
Description=mlbot monitor monthly drift
After=docker.service network-online.target
Requires=docker.service
OnFailure=mlbot-monitor-notify@%n.service

[Service]
Type=oneshot
ExecStart=/opt/quant-engine/deploy/systemd/mlbot-monitor-docker-run.sh monthly

[Install]
WantedBy=multi-user.target
